# custom/layers/meta-disco-scheduler/files/a53-app-sys-manager.service
[Unit]
Description=A53 Application System Manager

[Service]
ExecStart=/usr/bin/a53-app-sys-manager
Restart=always

[Install]
WantedBy=multi-user.target
